# templates/init.d/dynaTraceCollector.erb
#!/bin/sh
### Start/Stop dynaTrace Collector daemon
#
# Copyright 2004-2012 dynaTrace Software
# All rights reserved.
#
# Copy this script to /etc/init.d and use chkconfig for enabling/disabling
# the daemon. In order to run multiple Servers/Collectors on the same host,
# further adjustments may be necessary.
#
### BEGIN INIT INFO
# Provides: dynaTraceCollector
# Required-Start: $network
# Required-Stop: $network
# X-UnitedLinux-Should-Start:
# Default-Start: <%= @linux_service_start_runlevels %>
# Default-Stop: <%= @linux_service_stop_runlevels %>
# Description: Start/Stop dynaTrace Collector daemon
### END INIT INFO

# To use this file in your environment, adjusting the following variables:
# - DT_HOME     installation directory of dynaTrace software
# - DT_OPTARGS  optional arguments (e.g. -vm, -instance, -server, etc.)
# - DT_INSTANCE instance name if -instance argument is specified
DT_HOME=<%= @installer_prefix_dir %>/dynatrace
DT_OPTARGS="-listen <%= @agent_port %> -server <%= @server_hostname %>:<%= @server_port %><% if $jvm_xms { %> -Xms<%= @jvm_xms %><% } %><% if $jvm_xmx { %> -Xmx<%= @jvm_xmx %><% } %><% if $jvm_perm_size { %> -XX:PermSize=<%= @jvm_perm_size %><% } %><% if $jvm_max_perm_size { %> -XX:MaxPermSize=<%= @jvm_max_perm_size %><% } %>"
# When using an instance name, ensure that the -instance argument from
# DT_OPTARGS is in sync with the following DT_INSTANCE variable:
# (e.g.: "-instance InstName" requires "DT_INSTANCE=InstName")
DT_INSTANCE=

DT_PRODUCT=Collector
DT_BINARY=dtcollector

DT_PIDFILE=/var/run/$DT_BINARY.pid
if [ -n "$DT_INSTANCE" ]; then
  DT_PIDFILE=/var/run/$DT_BINARY.$DT_INSTANCE.pid
fi
if [ ! -x "$DT_HOME/$DT_BINARY" ]; then
  echo "File not found or not executable: $DT_HOME/$DT_BINARY"
  exit 1
fi
if [ -f "$DT_PIDFILE" ]; then
  PROCESSPID=$(cat $DT_PIDFILE)
fi
if [ -z "$PROCESSPID" ]; then
  GREPARG=`echo $DT_BINARY | cut -c1-8`
  if [ -z "$DT_INSTANCE" ]; then
    PROCESSPID=`ps -e | grep $GREPARG | awk '{{print $1}}'`
  else
    PROCESSPID=`ps -ef | grep $GREPARG | grep $DT_INSTANCE | awk '{{print $2}}'`
  fi
fi

case "$1" in
'start')
  if [ -z "$PROCESSPID" ]; then
    if [ -z "$DT_INSTANCE" ]; then
      echo ${DT_HOME}/${DT_BINARY} -bg $DT_OPTARGS
      su - dynatrace -c "${DT_HOME}/${DT_BINARY} -bg $DT_OPTARGS"
    else
      echo ${DT_HOME}/${DT_BINARY} -instance $DT_INSTANCE -bg $DT_OPTARGS
      su - dynatrace -c "${DT_HOME}/${DT_BINARY} -instance $DT_INSTANCE -bg $DT_OPTARGS"
    fi    
  else
    echo "dynaTrace $DT_PRODUCT daemon is already running!"
  fi
  ;;
'stop')
  if [ -n "$PROCESSPID" ]; then
    echo "Terminating dynaTrace $DT_PRODUCT process $PROCESSPID"
    kill -15 $PROCESSPID
  fi
  ;;
'version')
  $DT_HOME/$DT_BINARY -v
  ;;
'status')
  if [ -n "$PROCESSPID" ]; then
    echo "dynaTrace $DT_PRODUCT daemon is running:"
    ps -ef | grep $DT_BINARY | grep $PROCESSPID
      # | grep -v org.eclipse.equinox
  else
    echo "dynaTrace $DT_PRODUCT daemon not running."
    exit 1
  fi
  ;;
*)
  echo "usage: $0 {start|stop|status|version}"
  ;;
esac

exit 0
