# templates/init.d/dynaTraceBackendServer.erb
#!/bin/sh
### Start/Stop Dynatrace Server daemon
#
# Copyright 2004-2016 Dynatrace Software
# All rights reserved.
#
# Copy this script to /etc/init.d and use chkconfig for enabling/disabling
# the daemon. Depending on the target environment adjustments may be necessary.
#
### BEGIN INIT INFO
# Provides: dynaTraceServer
# Required-Start: $network
# Required-Stop: $network
# X-UnitedLinux-Should-Start:
# Default-Start: <%= @params["linux_service_start_runlevels"] %>
# Default-Stop: <%= @params["linux_service_stop_runlevels"] %>
# Description: Start/Stop Dynatrace Server daemon
### END INIT INFO

# To use this file in your environment, adjusting the following variables:
# - DT_HOME     installation directory of Dynatrace software
# - DT_OPTARGS  optional arguments (e.g. -vm, -server, etc.)
# - DT_RUNASUSER  the user to run the Dynatrace process

DT_HOME=<%= @params["installer_prefix_dir"] %>/dynatrace
DT_OPTARGS=

DT_PRODUCT=Server
DT_BINARY=dtserver

DT_RUNASUSER=<%= @params["user"] %>

# How long to wait for shutdown
DT_SHUTDOWN_WAIT=60

DT_PIDFILE_DIR=<%= @dynatrace::pid_file_directory %>

DT_PIDFILE=${DT_PIDFILE_DIR}/${DT_BINARY}63.pid

if [ ! -x "$DT_HOME/${DT_BINARY}" ]; then
  echo "File not found or not executable: ${DT_HOME}/${DT_BINARY}"
  exit
fi
if [ -f "${DT_PIDFILE}" ]; then
  PROCESSPID=$(cat ${DT_PIDFILE})
fi
if [ -z "${PROCESSPID}" ]; then
  RPID=`ps -A -o pid,comm | grep ${DT_BINARY} | awk '{ printf "%s", $1 }'`
  #RPID=$(ps -C ${DT_BINARY} -o pid --no-headers | awk '{printf "%s|", $0}' | awk '{sub(/\|$/,""); printf $0}')
  if [ ! -z "${RPID}" ]; then
    PROCESSPID=$(ps -ef | grep ${DT_BINARY} | awk '{{print $2}}' | grep ${RPID})
  fi
fi

case "$1" in
'start')
  if [ -z "${PROCESSPID}" ]; then
    if [ -z "${DT_RUNASUSER}" ]; then
      echo "${DT_HOME}"/${DT_BINARY} -bg ${DT_OPTARGS}
      "${DT_HOME}"/${DT_BINARY} -bg ${DT_OPTARGS}
      echo $! >${DT_PIDFILE}
    else
      echo "${DT_HOME}"/${DT_BINARY} -bg ${DT_OPTARGS}
      su -c "${DT_HOME}/${DT_BINARY} -bg ${DT_OPTARGS}" "${DT_RUNASUSER}"
      echo $! >${DT_PIDFILE}
    fi
  else
    echo "Dynatrace ${DT_PRODUCT} daemon is already running!"
  fi
  ;;
'stop')
  if [ -n "${PROCESSPID}" ]; then
    echo "Terminating Dynatrace $DT_PRODUCT process ${PROCESSPID}"
    kill -15 ${PROCESSPID}

    COUNT=0;
    while [ `ps -A -o pid | grep -c ${PROCESSPID}` -gt 0 ] && [ "${COUNT}" -lt "${DT_SHUTDOWN_WAIT}" ] # `ps --pid ${PROCESSPID} | grep -c ${PROCESSPID}` -ne 0]
    do
      echo "Waiting for Dynatrace $DT_PRODUCT (${PROCESSPID}) to finish shutdown";
      sleep 1
      COUNT=`expr ${COUNT} + 1`
    done

    if [ "${COUNT}" -gt "${DT_SHUTDOWN_WAIT}" ]; then
      echo "Killing Dynatrace $DT_PRODUCT (${PROCESSPID}) because the shutdown lasted longer than ${DT_SHUTDOWN_WAIT} seconds"
      kill -9 ${PROCESSPID}
    fi

    if [ -f "${DT_PIDFILE}" ]; then
      rm -f ${DT_PIDFILE}
    fi
  fi
  ;;
'version')
  "${DT_HOME}"/${DT_BINARY} -v
  ;;
'status')
  if [ -n "${PROCESSPID}" ]; then
    echo "Dynatrace ${DT_PRODUCT} daemon is running:"
    ps -ef | grep ${DT_BINARY} | grep ${PROCESSPID}
  else
    echo "Dynatrace ${DT_PRODUCT} daemon not running."
    exit 1
  fi
  ;;
*)
  echo "usage: $0 {start|stop|status|version}"
  ;;
esac

exit 0
